# debian/postinst
#!/bin/bash
set -e

if [ "$1" = "configure" ]; then
    # Install TensorFlow via pip if not available
    if ! python3 -c "import tensorflow" 2>/dev/null; then
        echo "Installing TensorFlow via pip..."
        python3 -m pip install --break-system-packages tensorflow>=2.10.0 || true
    fi
    
    systemctl daemon-reload
    echo "SSHGuard installed. Enable with: systemctl enable sshguard"
fi

#DEBHELPER#
exit 0
